# etc/systemd/system/VMAX_AF.socket
[Unit]
Description=VMAX Docker Volume Plugin

[Socket]
ListenStream=/run/docker/plugins/VMAX_AF.sock

[Install]
WantedBy=sockets.target